# push_up_tracker.service
[Unit]
Description=Push Up Tracker Web Application
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/push_up_tracker
ExecStart=/opt/push_up_tracker/push_up_tracker
Restart=always
RestartSec=10

# Load environment from .env file
EnvironmentFile=/opt/push_up_tracker/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/push_up_tracker

# Additional security hardening
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateDevices=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target